# rootfs/etc/services.d/backup/run
#!/command/with-contenv bash
# Periodic backup service using restic

# Source s6 container environment (ensures vars are available)
for f in /run/s6/container_environment/RESTIC_*; do
  [ -f "$f" ] && export "$(basename $f)"="$(cat $f)"
done

CONFIG_FILE="/etc/openclaw/backup.yaml"

# Skip if ENABLE_SPACES is not true
if [ "${ENABLE_SPACES:-false}" != "true" ]; then
  echo "[backup] Disabled (set ENABLE_SPACES=true for persistence)"
  s6-svc -Od .
  exit 0
fi

# Skip if persistence not configured
if [ -z "$RESTIC_REPOSITORY" ] || [ -z "$RESTIC_PASSWORD" ]; then
  echo "[backup] ENABLE_SPACES=true but RESTIC_REPOSITORY or RESTIC_PASSWORD not set"
  s6-svc -Od .
  exit 0
fi

# Read backup interval from config (default 30 seconds)
BACKUP_INTERVAL=$(yq -r '.backup_interval_seconds // 30' "$CONFIG_FILE")
echo "[backup] Starting backup loop (every ${BACKUP_INTERVAL}s)..."

while true; do
  sleep "$BACKUP_INTERVAL"
  /usr/local/bin/restic-backup
done
